# batsignal/batsignal.service.d/options.conf
[Service]
ExecStart=
ExecStart=/run/current-system/sw/bin/batsignal -c 10 -w 30 -f 97
